# deploy/systemd/vllm-gpt-oss-20b.service
[Unit]
Description=vLLM GPT-OSS 20B server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hydralisk
Group=hydralisk
WorkingDirectory=/opt/hydralisk
EnvironmentFile=/etc/hydralisk/hydralisk.env
Environment=HF_HOME=/var/lib/hydralisk/huggingface
Environment=PATH=/opt/hydralisk/.venv/bin:/usr/local/cuda/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=VLLM_USAGE_SOURCE=production
ExecStart=/opt/hydralisk/.venv/bin/vllm serve openai/gpt-oss-20b --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/var/lib/hydralisk /opt/hydralisk

[Install]
WantedBy=multi-user.target
